# quartz-top/.config/zsh/rc.local
#!/usr/bin/env sh

alias dwn="hyprctl keyword monitor HDMI-A-1, addreserved, 60,0,0,0"
